# src/conf/supervisor.conf
[supervisord]
nodaemon=true

[program:mailcatcher]
command=/bin/bash -c "mailcatcher --foreground --ip=0.0.0.0"
